# host/pal-watch.service
[Unit]
Description=Check Palworld host health and send PushPlus alerts

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pal-watch
